# agent/agent.service
[Unit]
Description=SNMP Agent Service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/snmp-agent
ExecStart=/usr/bin/python3 /home/ubuntu/snmp-agent/service_runner.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
